# pootle/typo3.cron
# Cron jobs for Pootle and TYPO3

# 3rd party extensions
15 2     * * *     pootle   /home/pootle/scripts/bin/build-language-pack 2>&1 >/tmp/build-all.log

# Backup
0 3      * * *     pootle   /home/pootle/scripts/bin/backup-database >/dev/null 2>&1
